# systemd/hindsight-embed.service
[Unit]
Description=Hindsight profile embedding job

[Service]
Type=oneshot
# ExecStart must be the absolute path of your installed `hindsight` binary.
# The default below assumes `~/.local/bin/hindsight`; edit it to match your
# install location. This is a drain-and-exit batch run: it assembles synthetic
# profiles and embeds them into the vector store, then exits. It is deliberately
# NOT the capture daemon (a long deferring embed would fight the daemon's
# 15-minute idle self-terminate), and it is triggered by hindsight-embed.timer.
ExecStart=%h/.local/bin/hindsight embed
